# crontab
MAILTO=root
# run dummy script every 1 minutes
* * * * */1 /usr/bin/php /scripts/dummy.php
